# deploy/systemd/novaadapt-runtime.service
[Unit]
Description=NovaAdapt Built-in Runtime Endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=novaadapt
Group=novaadapt
WorkingDirectory=/opt/novaadapt
EnvironmentFile=-/etc/novaadapt/runtime.env
ExecStart=/opt/novaadapt/deploy/systemd/start-runtime.sh
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/novaadapt /run/novaadapt
RuntimeDirectory=novaadapt
RuntimeDirectoryMode=0750
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
